# launchd/com.cnlangplugin.autopatch.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__LABEL__</string>

    <!-- Run the check when Claude.app's Info.plist changes (a Squirrel update
         replaces the whole bundle, recreating this file) and once at login. -->
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__CLAUDE_ZH_BIN__</string>
        <string>autopatch</string>
        <string>run</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>__APP_INFO_PLIST__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Coalesce rapid file-change bursts during an update. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>__LOG__</string>
    <key>StandardErrorPath</key>
    <string>__LOG__</string>
</dict>
</plist>
